# deploy/oak-drive-controller.service
# Systemd service for the Pi drive controller
# NOTE: The install.sh script generates this file with correct paths.
#       This is a reference template only — do not install directly.
#
# Manual install (adjust paths for your user):
#   sudo cp deploy/oak-drive-controller.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now oak-drive-controller

[Unit]
Description=OAK-Drive-Sync Pi Drive Controller
After=network-online.target mosquitto.service
Wants=network-online.target
Requires=mosquitto.service

[Service]
Type=exec
User=pi
Group=pi
WorkingDirectory=/home/pi/OAK-Drive-Sync
ExecStart=/home/pi/.local/bin/uv run oak-pi run --config config/drive_config.yaml
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# GPIO access
SupplementaryGroups=gpio

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/home/pi/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
